# systemd/system/stock-market-words.service
[Unit]
Description=Stock Market Words Data Pipeline
Documentation=https://github.com/aallbright/stock-market-words
After=network-online.target pia-rotate.service
Wants=network-online.target

[Service]
Type=oneshot
User=smw
Group=smw
WorkingDirectory=/opt/stock-market-words/python3

# Environment setup
Environment="PATH=/opt/stock-market-words/python3/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Run the pipeline (ticker-cli entry point installed by pip install -e .)
ExecStart=/opt/stock-market-words/python3/venv/bin/ticker-cli run-all

# Git commit and push changes (if configured)
ExecStartPost=/opt/stock-market-words/systemd/scripts/git-commit-and-push.sh

# Back up database to S3 (optional — requires aws CLI + credentials)
# Uncomment after completing AWS setup per docs/design/20260417_015817_UTC_s3_database_backup.md
# ExecStartPost=/opt/stock-market-words/systemd/scripts/backup-db-s3.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stock-market-words

# Resource limits
TimeoutStartSec=3600
MemoryMax=4G
CPUQuota=200%

# Restart policy
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
