# systemd/system/stock-market-words.timer
[Unit]
Description=Stock Market Words Data Pipeline Timer (Daily at 6PM EST)
Documentation=https://github.com/aallbright/stock-market-words
Requires=stock-market-words.service

[Timer]
# Run at 23:00 UTC on weekdays only (Mon-Fri)
# - During EST (Nov-Mar): 23:00 UTC = 6:00 PM ET — 2h after market close
# - During EDT (Mar-Nov): 23:00 UTC = 7:00 PM ET — 3h after market close
# US market holidays are harmless: pipeline is idempotent, produces no new data
OnCalendar=Mon..Fri *-*-* 23:00:00

# If system was down, run on next boot (but wait 5 minutes)
Persistent=true
OnBootSec=5min

# Randomize start time by up to 5 minutes to avoid thundering herd
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
